# recipes-support/swupdate/swupdate-2019.04/swupdate-progress.service
[Unit]
Description=swupdate progress service
Requires=swupdate.service

[Service]
ExecStart=/usr/bin/swupdate-progress -r -w
